# deploy/desktop/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

# =============================================================================
# Virtual desktop (TigerVNC :99 + XFCE + noVNC on 6080→5900)
# =============================================================================
[program:desktop]
command=bash /app/desktop/desktop.sh
user=unityuser
environment=HOME="/Droid",DISPLAY=":99",XDG_RUNTIME_DIR="/tmp/runtime-unityuser",DBUS_SESSION_BUS_ADDRESS="autolaunch:"
autostart=true
autorestart=true
priority=100
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/desktop.log
stderr_logfile=/var/log/supervisor/desktop-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# =============================================================================
# Virtual audio devices (self-host Meet / voice plumbing)
# =============================================================================
[program:device]
command=bash /app/desktop/device.sh
autostart=true
autorestart=true
priority=200
stdout_logfile=/var/log/supervisor/device.log
stderr_logfile=/var/log/supervisor/device-err.log

# =============================================================================
# Agent Service
# =============================================================================
[program:agent-service]
command=npx ts-node src/index.ts
directory=/app/agent-service
user=unityuser
environment=NODE_ENV="production",DISPLAY=":99",HOME="/Droid",DROID_LOCAL_ROOT="/Droid/Local",MAGNITUDE_LOG_DIR="/var/log/magnitude",MAGNITUDE_DEBUG="true",UNIFY_KEY="%(ENV_UNIFY_KEY)s",DROID_GATEWAY_URL="%(ENV_DROID_GATEWAY_URL)s",ORCHESTRA_URL="%(ENV_ORCHESTRA_URL)s",DROID_COMMS_URL="%(ENV_DROID_COMMS_URL)s",UNIFY_MODEL="%(ENV_UNIFY_MODEL)s"
autostart=true
autorestart=true
priority=300
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/agent-service.log
stderr_logfile=/var/log/supervisor/agent-service-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# =============================================================================
# SFTP file sync (unityuser chroot at /Droid)
# =============================================================================
[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
priority=400
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd-err.log

[group:desktop]
programs=desktop
priority=10

[group:services]
programs=device,agent-service,sshd
priority=20
