# desktop/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:desktop]
command=bash /app/desktop/desktop.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/desktop.log
stderr_logfile=/var/log/supervisor/desktop-err.log
priority=100

[program:device]
command=bash /app/desktop/device.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/device.log
stderr_logfile=/var/log/supervisor/device-err.log
priority=200

[program:agent-service]
command=npx ts-node src/index.ts
directory=/app/agent-service
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/agent-service.log
stderr_logfile=/var/log/supervisor/agent-service-err.log
priority=300
